# crontab
PYTHONPATH=/usr/local/lib/python2.7/site-packages
0 */6 * * * . /tmp/.env; python /usr/local/bin/github-backup $GITHUB_USERNAME -o /storage -t $GITHUB_TOKEN --repositories >> /var/log/cron.log 2>&1
